# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Store Service (Django)
[program:store]
command=sh -c "python manage.py collectstatic --no-input && gunicorn --bind 0.0.0.0:8001 store_service.wsgi:application"
directory=/app/store
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Cargo Service (FastAPI)
[program:cargo]
command=uvicorn app.main:app --host 0.0.0.0 --port 8002
directory=/app/cargo
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Hotels Service (Django)
[program:hotels]
command=sh -c "python manage.py collectstatic --no-input && gunicorn --bind 0.0.0.0:8003 hotels_service.wsgi:application"
directory=/app/hotels
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Language Academy (FastAPI)
[program:language_academy]
command=uvicorn app.main:app --host 0.0.0.0 --port 8004
directory=/app/language_academy
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# AI Assistant (FastAPI)
[program:ai_assistant]
command=uvicorn app.main:app --host 0.0.0.0 --port 8005
directory=/app/ai_assistant
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Frontend (FastAPI + Jinja2)
[program:frontend]
command=uvicorn app:app --host 0.0.0.0 --port 3000
directory=/app/frontend
environment=PORT="3000",STORE_BACKEND_URL="http://127.0.0.1:8001/api/products/",STORE_PAYMENTS_URL="http://127.0.0.1:8001/payments/process/",AI_ASSISTANT_URL="http://127.0.0.1:8005",FRONTEND_IYOBO_URL="/chat"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
